# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:api]
command=/bin/sh -c 'python -m uvicorn apps.api.main:app --host "${API_HOST:-0.0.0.0}" --port "${API_PORT:-8000}"'
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=/bin/sh -c 'HOSTNAME="${WEB_HOST:-0.0.0.0}" PORT="${WEB_PORT:-3000}" node server.js'
directory=/app/apps/web
autostart=true
autorestart=true
priority=20
startsecs=5
environment=NEXT_PUBLIC_API_BASE_URL="%(ENV_NEXT_PUBLIC_API_BASE_URL)s",NEXT_PUBLIC_API_PREFIX="%(ENV_NEXT_PUBLIC_API_PREFIX)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001
